# DAI/Practica_9/web/gunicorn.ini
[program:gunicorn]
command=/usr/local/bin/gunicorn sitio_web.wsgi:application  --bind 0.0.0.0:8000
directory=/web
user=root
autostart=true
autorestart=true
redirect_stderr=true
